# etc/crontab
# m h  dom mon dow   command
CWD=$(pwd)
0 5  * * *  /home/touchpointmed/nginx-hapi-redis/cron_job.sh
