# etc/crontab
# m h  dom mon dow   command
0 5  * * *  /home/touchpointmed/sandbox/cron_job.sh
