# etc/crontab
# m h  dom mon dow   command
0 5  * * *  $CWD/cron_job.sh
